# container-files/etc/supervisor.d/filebeat.conf
[program:filebeat]
command = /usr/bin/filebeat -e
autorestart = true