# container-files/etc/supervisor.d/varnishd.conf
[program:varnishd]
command = /usr/sbin/varnishd -F -u varnish -f /etc/varnish/default.vcl -s malloc,64m -p default_ttl=3600 -p default_grace=3600
autorestart = true