# container-files/etc/supervisor.d/varnishd.conf
[program:varnishd]
command = /usr/sbin/varnishd -F -f %(ENV_VCL_CONFIG)s -s malloc,%(ENV_CACHE_SIZE)s %(ENV_VARNISHD_PARAMS)s
autorestart = true